# systemd/wattometer-server.service
[Unit]
Description=Wattometer data server

[Service]
ExecStart=/home/tami/dev/power-meter/server/socat.sh

[Install]
WantedBy=default.target